# root/etc/systemd/system/tomcat8@.service.d/webtop.conf
[Service]
ExecStartPre=-/usr/bin/find /var/lib/tomcats/webtop/work/Catalina/localhost/webtop/ -delete
